# src/influxdb/supervisord.conf
[supervisord]
nodaemon=true
[program:influxdb]
command = sh -c "influxd"
autorestart = true

[program:telegraf]
command = sh -c "telegraf --config /telegraf.conf"
autorestart = true
