# src/nginx/supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command = sh -c "/usr/sbin/sshd -D"
autorestart = true

[program:nginx]
command = sh -c "nginx -g 'daemon off;'"
autorestart = true
